# rootfs/etc/s6-overlay/s6-rc.d/init-ssh/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# PowerHaus SSH Initialization
# Adapted from hassio-addons/addon-ssh
# Generates host keys, configures sshd, sets up user, generates ttyd credential
# ==============================================================================

declare username

bashio::log.info "Initializing SSH..."

# --- Generate host keys if they don't exist ---
if ! bashio::fs.file_exists "/data/ssh_host_rsa_key"; then
    bashio::log.info "Generating RSA host key..."
    ssh-keygen -t rsa -b 4096 -f /data/ssh_host_rsa_key -N "" \
        || bashio::exit.nok "Failed to generate RSA host key"
fi

if ! bashio::fs.file_exists "/data/ssh_host_ed25519_key"; then
    bashio::log.info "Generating ED25519 host key..."
    ssh-keygen -t ed25519 -f /data/ssh_host_ed25519_key -N "" \
        || bashio::exit.nok "Failed to generate ED25519 host key"
fi

# --- Read configuration ---
username=$(bashio::config 'ssh.username')

# --- Create/update user ---
if id -u "${username}" > /dev/null 2>&1; then
    bashio::log.info "User '${username}' already exists, updating..."
else
    bashio::log.info "Creating user '${username}'..."
fi

# Create user with home directory
adduser -D -s /bin/bash "${username}" 2>/dev/null || true

# --- Always disable password auth (Studio manages keys) ---
sed -i "s/PasswordAuthentication.*/PasswordAuthentication no/" /etc/ssh/sshd_config
bashio::log.info "Password authentication disabled (Studio-managed SSH keys only)."

# --- Setup authorized keys directory ---
# Keys will be populated dynamically by the config sync task in app.py
mkdir -p "/home/${username}/.ssh"
touch "/home/${username}/.ssh/authorized_keys"
chmod 600 "/home/${username}/.ssh/authorized_keys"

# Also add any locally-configured keys from add-on options
if bashio::config.has_value 'ssh.authorized_keys'; then
    bashio::log.info "Adding locally-configured authorized keys..."
    while IFS= read -r key; do
        echo "${key}" >> "/home/${username}/.ssh/authorized_keys"
    done <<< "$(bashio::config 'ssh.authorized_keys')"
fi

chown -R "${username}:${username}" "/home/${username}/.ssh"

# --- Configure sshd ---
# Set host keys
sed -i "s|#HostKey /data/ssh_host_rsa_key|HostKey /data/ssh_host_rsa_key|" /etc/ssh/sshd_config
sed -i "s|#HostKey /data/ssh_host_ed25519_key|HostKey /data/ssh_host_ed25519_key|" /etc/ssh/sshd_config

# Set allowed user
sed -i "s/AllowUsers .*/AllowUsers ${username}/" /etc/ssh/sshd_config

# Configure TCP forwarding
if bashio::config.true 'ssh.allow_tcp_forwarding'; then
    sed -i "s/AllowTcpForwarding.*/AllowTcpForwarding yes/" /etc/ssh/sshd_config
else
    sed -i "s/AllowTcpForwarding.*/AllowTcpForwarding no/" /etc/ssh/sshd_config
fi

# Configure SFTP
if bashio::config.true 'ssh.sftp'; then
    sed -i "s|#Subsystem sftp|Subsystem sftp|" /etc/ssh/sshd_config
fi

# --- Generate ttyd internal credential ---
if ! bashio::fs.file_exists "/data/ttyd_credential"; then
    TTYD_PASS=$(head -c 32 /dev/urandom | base64 | tr -dc 'a-zA-Z0-9' | head -c 32)
    echo "${TTYD_PASS}" > /data/ttyd_credential
    chmod 600 /data/ttyd_credential
    bashio::log.info "Generated internal ttyd credential."
fi

# --- Setup user environment ---
# Create persistent data directory for user
mkdir -p /data/user_home

# Link useful Home Assistant directories to user home
for dir in config addons share backup media ssl; do
    if [ -d "/${dir}" ]; then
        ln -sf "/${dir}" "/home/${username}/${dir}" 2>/dev/null || true
    fi
done

# Persist shell history
if ! bashio::fs.file_exists "/data/user_home/.bash_history"; then
    touch /data/user_home/.bash_history
fi
ln -sf /data/user_home/.bash_history "/home/${username}/.bash_history"

# Set SUPERVISOR_TOKEN in user environment
echo "SUPERVISOR_TOKEN=${SUPERVISOR_TOKEN}" > "/home/${username}/.ssh/environment" 2>/dev/null || true
chmod 600 "/home/${username}/.ssh/environment" 2>/dev/null || true

# Ensure correct ownership
chown -R "${username}:${username}" "/home/${username}"

bashio::log.info "SSH initialization complete."
