# rootfs/etc/s6-overlay/s6-rc.d/powerhaus-app/finish
#!/usr/bin/execlineb -S1
# ==============================================================================
# PowerHaus App Finish Script
# Handles unexpected exits
# ==============================================================================

if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }
s6-svc -O .
