# powerhaus/rootfs/etc/s6-overlay/s6-rc.d/powerhaus-app/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# PowerHaus App Service
# Runs the main PowerHaus aiohttp application
# ==============================================================================

bashio::log.info "Starting PowerHaus application..."

exec python3 /app/app.py
