# rootfs/etc/s6-overlay/s6-rc.d/sshd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# PowerHaus SSH Daemon
# Runs the OpenSSH server
# ==============================================================================

bashio::log.info "Starting SSH daemon..."

exec /usr/sbin/sshd -D -e
